# nic-reminder.service
[Unit]
Description=NIC Call Center Reminder Service
Documentation=https://github.com/Vkdel001/Callcenter
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/nic-callcenter
ExecStart=/usr/bin/node /var/www/nic-callcenter/backend-reminder-service.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nic-reminder

# Environment variables
Environment=NODE_ENV=production
EnvironmentFile=/var/www/nic-callcenter/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/www/nic-callcenter

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target